# debian/darkplug/DEBIAN/postinst
#!/bin/sh

set -e
systemctl daemon-reload
systemctl enable darkplug.service
systemctl start darkplug.service
