# debian/darkplug/DEBIAN/prerm
#!/bin/sh

# set -e # allowed to fail
systemctl stop darkplug.service
systemctl disable darkplug.service
